# crontab
* * * * * /bin/bash /SmartUP/scripts/trigger-scraper.sh >> /var/log/cron.log 2>&1
